# launchagents/dev.capslock-blink.poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>dev.capslock-blink.poll</string>
  <key>ProgramArguments</key>
  <array><string>/bin/bash</string><string>__DIR__/poll.sh</string></array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardErrorPath</key><string>/tmp/capsled-poll.err.log</string>
  <key>StandardOutPath</key><string>/tmp/capsled-poll.out.log</string>
</dict>
</plist>
